# docker/production/php/cli/supervisor/cron.conf
[program:cron-worker]
process_name=%(program_name)s_%(process_num)02d
command=cron -f
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/cron-worker.log
